# etc/init/remount.conf
# Executes remount on user login.
# start-user-session is a signal emmitted by session_manager

start on start-user-session

script
        sleep 2
        exec /opt/bin/remount
end script
