# etc/init/remount.conf
# Remount every mount point after user has logged in

start on stopped login 

script
        exec /opt/bin/remount
end script
